# docker/supervisor.conf
[program:node]
command=/usr/bin/node /var/www/ivr-server.js
numprocs=1
directory=/var/www/
stdout_logfile=/var/log/ivr/ivr.log
redirect_stderr=true
autostart=true
autorestart=true